# systemd/github_rate_limit_exporter.service
[Unit]
Description=GitHub Rate Limit Exporter
Documentation=https://github.com/l13t/github_rate_limit_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=prometheus
Group=prometheus
ExecStart=/usr/local/bin/github_rate_limit_exporter -config /etc/github_rate_limit_exporter/config.yaml
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/github_rate_limit_exporter
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=8192
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github_rate_limit_exporter

[Install]
WantedBy=multi-user.target
